# com.fs2-menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fs2-menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>VENV_PYTHON</string>
        <string>MAIN_PY</string>
    </array>
    <key>WorkingDirectory</key>
    <string>WORKING_DIR</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/fs2-menubar.err.log</string>
</dict>
</plist>
